# GroupBuya/GroupBuya - Grab BrandReward exclusive + favor hk.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2021-09-08T17:02:22.6289106</Date>
    <Author>admin</Author>
    <URI>\GroupBuya - Grab BrandReward exclusive + favor hk</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT12H</ExecutionTimeLimit>
    <MultipleInstancesPolicy>StopExisting</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2021-09-08T17:01:25</StartBoundary>
      <ScheduleByWeek>
        <WeeksInterval>1</WeeksInterval>
        <DaysOfWeek>
          <Monday />
          <Wednesday />
          <Friday />
        </DaysOfWeek>
      </ScheduleByWeek>
    </CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>C:\project\yii2\GroupBuya\yii.bat</Command>
      <Arguments>grabpromocode brandreward exclusive</Arguments>
      <WorkingDirectory>C:\project\yii2\GroupBuya\</WorkingDirectory>
    </Exec>
    <Exec>
      <Command>C:\project\yii2\GroupBuya\yii.bat</Command>
      <Arguments>grabpromocode brandreward favorhk</Arguments>
      <WorkingDirectory>C:\project\yii2\GroupBuya\</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
